# module1/task2/etc_backup_7zip/etc_backup/systemd/system/multi-user.target.wants/lxd-installer.socket
[Unit]
Description=Helper to install lxd snap on demand

[Socket]
ListenStream=/run/lxd-installer.socket
SocketUser=root
SocketGroup=lxd
SocketMode=0660
Accept=true

[Install]
WantedBy=multi-user.target
